# supervisord.conf

[program:taskserver]
command=task-execution-server/bin/ga4gh-taskserver -db var/ga4gh_tasks.db -storage var/storage

[program:worker]
command=task-execution-server/bin/ga4gh-worker
[supervisord]


[unix_http_server]
file=./supervisor.sock 

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://./supervisor.sock ; use a unix:// URL  for a unix socket
