# supervisord.conf

[program:taskserver]
command=task-execution-server/bin/ga4gh-taskserver -db var/ga4gh_tasks.db -storage var/storage
stdout_logfile=var/taskserver.stdout
stderr_logfile=var/taskserver.stderr

[program:worker]
command=sleep 5 && task-execution-server/bin/ga4gh-worker -storage var/storage
stdout_logfile=var/worker.stdout
stderr_logfile=var/worker.stderr

[supervisord]


[unix_http_server]
file=./supervisor.sock 

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://./supervisor.sock ; use a unix:// URL  for a unix socket
